# scripts/configs/amplipi.service
[Unit]
Description=AmpliPi Audio Controller
After=network.target sound.target
Wants=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi
ExecStart=/usr/local/bin/amplipi --addr :80
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=amplipi

# Allow binding to privileged ports (< 1024) as non-root user
AmbientCapabilities=CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target
